# deploy/nopmsbot-v2.service
[Unit]
Description=NoPMsBot — Telegram support-desk bot (v2.3)
After=network-online.target
Wants=network-online.target
# Cap the crash-restart loop: after 5 failures within 60s, systemd stops
# retrying instead of looping forever on a broken env/config (M10,
# docs/AUDIT-2026-07-10.md). `systemctl reset-failed nopmsbot-v2` re-arms it.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=nopmsbot
Group=nopmsbot
WorkingDirectory=/opt/nopmsbot-v2
EnvironmentFile=/etc/nopmsbot-v2.env
ExecStart=/opt/nopmsbot-v2/.venv/bin/python bot.py
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/nopmsbot-v2

[Install]
WantedBy=multi-user.target
